# supervisor/initialize.ini
[program:initialize]
command=/usr/bin/start-all.sh restart
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0