# supervisor/kerberos.ini
[program:krb5kdc]
command = /usr/sbin/krb5kdc -n
autostart=false
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:kaasdmind]
command = /usr/sbin/kadmind -nofork
autostart=false 
autorestart=false
redirect_stderr = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0